# ATA-1/etc/crontabs/ntpc

#Time Sync
0 5 * * * /etc/rc.d/init.d/ntpc > /dev/null 2>&1


